# systemd/nginx-proxy.socket
[Socket]
ListenStream=/opt/containers/run/*

[Install]
WantedBy=sockets.target
