# systemd/nginx-proxy.service
[Unit]
Requires=nginx-docker.service
After=nginx-docker.service

[Service]
ExecStart=/opt/bin/launcher.sh
